# playground/supervisord.conf
[supervisord]
nodaemon=true
loglevel=trace
# logfile=/dev/null
# logfile_maxbytes=0

[program:cat]
command=cat

[unix_http_server]
chmod = 0777
chown= nobody:nogroup

[inet_http_server]
port = 0.0.0.0:9001
